# dnsctl.service
[Unit]
Description=dnsctl — Technitium DNS management UI
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/dnsctl
EnvironmentFile=/root/dnsctl/.env
ExecStart=/root/dnsctl/.venv/bin/python /root/dnsctl/main.py
Restart=on-failure
RestartSec=5

# Redirect stdout/stderr to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dnsctl

[Install]
WantedBy=multi-user.target
