# harp.service
[Unit]
Description=harp — Technitium DNS management UI
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/harp
EnvironmentFile=/root/harp/.env
ExecStart=/root/harp/.venv/bin/python /root/harp/main.py
Restart=on-failure
RestartSec=5

# Redirect stdout/stderr to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=harp

[Install]
WantedBy=multi-user.target
